# whatnot-collector/deploy/systemd/user/whatnot-redis.service
[Unit]
Description=Whatnot Runtime Redis
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/cybertechna/AethrixSystems_Portable/hjay9672-WN /whatnot-collector
ExecStartPre=/bin/bash -lc 'mkdir -p "/home/cybertechna/AethrixSystems_Portable/hjay9672-WN /whatnot-collector/data/runtime/redis"'
ExecStart=/bin/bash -lc 'exec /opt/gitlab/embedded/bin/redis-server --bind 127.0.0.1 --port 6379 --dir "/home/cybertechna/AethrixSystems_Portable/hjay9672-WN /whatnot-collector/data/runtime/redis" --dbfilename whatnot-runtime.rdb --appendonly yes --appendfilename whatnot-runtime.aof --daemonize no --protected-mode yes'
Restart=always
RestartSec=3

[Install]
WantedBy=default.target
