# whatnot-collector/deploy/systemd/user/whatnot-celery-beat.service
[Unit]
Description=Whatnot Celery Beat Scheduler
After=network.target whatnot-redis.service
Wants=whatnot-redis.service

[Service]
Type=simple
WorkingDirectory=/home/cybertechna/AethrixSystems_Portable/hjay9672-WN /whatnot-collector
Environment=REDIS_ENABLED=1
Environment=REDIS_URL=redis://127.0.0.1:6379/0
Environment=CELERY_BROKER_URL=redis://127.0.0.1:6379/0
Environment=CELERY_RESULT_BACKEND=redis://127.0.0.1:6379/0
ExecStartPre=/bin/bash -lc 'mkdir -p "/home/cybertechna/AethrixSystems_Portable/hjay9672-WN /whatnot-collector/data/runtime"'
ExecStart=/bin/bash -lc 'exec ".venv/bin/celery" -A app.workers.beat.celery_app beat --loglevel=INFO --schedule data/runtime/celerybeat-schedule'
Restart=always
RestartSec=5
StandardOutput=append:/home/cybertechna/AethrixSystems_Portable/hjay9672-WN /whatnot-collector/data/celery_beat.log
StandardError=append:/home/cybertechna/AethrixSystems_Portable/hjay9672-WN /whatnot-collector/data/celery_beat.log

[Install]
WantedBy=default.target
